# deploy/system/quantlab-gpu-control.service
[Unit]
Description=QuantLab GPU control endpoint
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=quantlab
Group=docker
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d -o quantlab -g docker -m 0755 /run/quantlab
ExecStart=/usr/bin/python3 /home/quantlab/quantlab-ai-capital/deploy/system/quantlab_gpu_control_server.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/run/quantlab /home/quantlab/quantlab-ai-capital/.env /home/quantlab/quantlab-runtime/models

[Install]
WantedBy=multi-user.target
