# deploy/systemd/quantlab-paper-trading.timer
[Unit]
Description=Run QuantLab Paper Trading Cycle Every Minute

[Timer]
OnCalendar=*:0/1
AccuracySec=10s
Persistent=true
Unit=quantlab-paper-trading.service

[Install]
WantedBy=timers.target
