# deploy/systemd/quantlab-night-mode.service
[Unit]
Description=QuantLabs night low-power mode
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/home/quantlab/quantlab-ai-capital/deploy/systemd/quantlab-power-mode.sh night
